# autonomy/nanoslot/install/nanoslot@.service
[Unit]
Description=Nanoslot meta handler
[Service]
ExecStart=/nanoslot/nanoslot.sh %i
User=robot

